# config/w84it.ini
[program:uwsgi]
command=uwsgi uwsgi.ini
stdout_logfile=/var/log/uwsgi.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/var/log/uwsgi_errors.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

